# system.d/iron-carrier.service
[Unit]
Description=Iron Carrier Daemon
After=network-online.target nss-lookup.target

[Service]
ExecStart=/usr/bin/iron-carrier -dv
Restart=on-failure
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
